# virtualkeyboard-mqtt.service
[Unit]
Description=MQTT Virtual Keyboard

[Service]
ExecStart=/usr/local/sbin/virtualkeyboard-mqtt

[Install]
WantedBy=multi-user.target